# dotfiles/eww/eww.service
[Unit]
Description=ElKowars wacky widgets
Documentation=https://elkowar.github.io/eww/
PartOf=graphical-session.target

[Service]
Type=exec
ExecStart=/usr/bin/eww -c %h/.config/eww daemon --no-daemonize
Restart=on-failure

[Install]
WantedBy=graphical-session.target
